# minecraft-update.timer
[Unit]
Description=Update Minecraft Server Binary

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timer.target
