# minecraft-restart.timer
[Unit]
Description=Restart Minecraft server

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timer.target
