# minecraft-update.timer
[Unit]
Description=Update Minecraft Server Binary

[Timer]
OnCalendar=*-*-* 4:00:00 PDT
Persistent=true

[Install]
WantedBy=timer.target
