# minecraft-snapshot@.service
[Unit]
Description=Snapshot Minecraft Worlds

[Service]
Type=oneshot
User=minecraft
Group=minecraft
ExecStart=/opt/minecraft/snapshot.sh %i
