# minecraft-reboot.timer
[Unit]
Description=Reboot server

[Timer]
OnCalendar=*-*-* 1:00:00 PDT
Unit=reboot.target

[Install]
WantedBy=timers.target
