# skel/mode644/etc/systemd/system/minecraft-watcher.path
[Path]
PathChanged=/opt/minecraft/server.jar

[Install]
WantedBy=multi-user.target
