# skel/mode644/etc/systemd/system/minecraft-update.timer
[Unit]
Description=Update Minecraft Server Binary

[Timer]
OnCalendar=00/1:00
Persistent=true

[Install]
WantedBy=timer.target
